# overlay/common/root/etc/services.d/apache2/finish
#!/usr/bin/execlineb -S0

echo "apache finished"

s6-svscanctl -t /var/run/s6/services