# overlay/common/root/etc/services.d/apache2/run
#!/usr/bin/with-contenv bash

echo "Starting Apache..."
exec apache2-foreground 
