# overlay/common/root/etc/cron.d/backup
0 * * * * www-data /overlay/bin/backup.sh backup > /proc/$(cat /var/run/crond.pid)/fd/1 2> /proc/$(cat /var/run/crond.pid)/fd/2
