# overlay/common/root/etc/cron.d/certbot
0 1 * * * root /overlay/bin/certbot.sh renew > /var/run/cron.stdout 2> /var/run/cron.stderr
