# overlay/common/root/etc/cron.d/certbot
0 1 * * * root /overlay/bin/certbot.sh renew > /proc/$(cat /var/run/crond.pid)/fd/1 2> /proc/$(cat /var/run/crond.pid)/fd/2
